# deploy/memi-ml-webhook.service
[Unit]
Description=Memi Mallorca deploy webhook
After=network.target

[Service]
User=root
WorkingDirectory=/opt/memi-ml
EnvironmentFile=/etc/memi-ml-webhook.env
ExecStart=/opt/memi-ml/.venv/bin/python deploy/webhook.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
